# supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:watch-recordings]
command=/usr/local/bin/watch-recordings.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/watch-recordings.log
stderr_logfile=/var/log/supervisor/watch-recordings-error.log

[program:cleanup-hls]
command=/usr/local/bin/cleanup-hls.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cleanup-hls.log
stderr_logfile=/var/log/supervisor/cleanup-hls-error.log
